# resources/dev-services/ghrunner.service
[Unit]
Description=GitHub Action Runner

[Service]
WorkingDirectory=/home/shane/actions-runner
ExecStart=/bin/bash run.sh
Restart=always

